# pi_server/plantower.cron
*/12 * * * * python /home/pi/PlantowerMonitor/main.py >> /home/pi/PlantowerMonitor/runtime.log 2>&1